# data/bootdrived.openrc
#!/sbin/openrc-run
# OpenRC service for the BootDrive backend (fallback for non-systemd installs).

name="bootdrived"
description="BootDrive USB mass-storage backend"
command="/usr/libexec/bootdrived"
command_background=true
pidfile="/run/bootdrived.pid"
command_user="root"

depend() {
    need dbus
    after usb-signaller
}

start_pre() {
    checkpath --directory --mode 0755 /run/bootdrived
}
